# templates/systemd-networkd-wait-online.service.j2
[Service]
ExecStart=/lib/systemd/systemd-networkd-wait-online -i {{ network_interface }}
